# pi-client/tv-guide.service
[Unit]
Description=Pi TV Guide Dashboard
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
User=ttumlin
WorkingDirectory=/home/ttumlin/pi-tv
EnvironmentFile=-/home/ttumlin/pi-tv/.env
ExecStart=/home/ttumlin/.local/bin/gunicorn -w 2 -b 0.0.0.0:5001 --timeout 3600 tv_guide:app
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tv-guide

[Install]
WantedBy=multi-user.target
