# pi-client/tv-player.service
[Unit]
Description=Pi TV Channel Player
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
User=ttumlin
Group=video
WorkingDirectory=/home/ttumlin/pi-tv
ExecStart=/usr/bin/python3 -u /home/ttumlin/pi-tv/tv_player.py
ExecStopPost=/bin/bash -c '/usr/bin/pkill -u ttumlin -f "mpv.*mpvsocket" || true; rm -f /tmp/mpvsocket'
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tv-player

# Give access to display and input devices
SupplementaryGroups=input video render

[Install]
WantedBy=multi-user.target
